# archive/2025-12-20-pre-consolidation/infrastructure/systemd/auto-free-space.timer
[Unit]
Description=Run auto-free-space daily and at boot

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
